# ops/systemd/clawpoint.service
[Unit]
Description=Clawpoint dashboard
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/clawpoint
EnvironmentFile=/opt/clawpoint/.env.production
ExecStart=/usr/bin/node /opt/clawpoint/server.mjs
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
