# ops/systemd/clawpoint.service
[Unit]
Description=Clawpoint dashboard
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/.local/share/clawpoint/current
EnvironmentFile=%h/.config/clawpoint/.env.production
ExecStart=/usr/bin/env node %h/.local/share/clawpoint/current/server.mjs
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
